# yar.service.example
[Unit]
Description=YAR XYJ Service
After=network.target

[Service]
Type=simple
User=netman
# Memory settings
MemoryHigh=8G
MemoryMax=12G

# Set the YAR installation directory (change this to match your installation path)
Environment="YAR_HOME=/home/netman/yar-xyj"

# Set Environment to your Python virtual environment
Environment="PATH=${YAR_HOME}/.venv/bin"
WorkingDirectory=${YAR_HOME}
ExecStart=${YAR_HOME}/.venv/bin/yar-server
# ExecStart=${YAR_HOME}/.venv/bin/yar-gunicorn

# Kill mode require ExecStart must be gunicorn or unvicorn main process
KillMode=process
ExecStop=/bin/kill -s TERM $MAINPID
TimeoutStopSec=60

Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
